# container/supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=/usr/sbin/nginx -g 'daemon off;'

[program:streamlit]
command=streamlit run /app/your_app.py --server.port=8501 --server.enableCORS=false